# daemon/clawd-runner.service
[Unit]
Description=Clawd Bot - Autonomous Claude Code Task Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user
ExecStart=/home/ec2-user/clawd-bot/daemon/clawd-runner.sh
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
Environment="HOME=/home/ec2-user"
Environment="PATH=/home/ec2-user/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/home/ec2-user/.clawd-env

[Install]
WantedBy=multi-user.target
